# archive/display_images.service
[Unit]
Description=Run display.py at user login
After=graphical.target
Wants=graphical.target

[Service]
Type=simple
WorkingDirectory=/home/multiplane/MultiPlaneDisplay
ExecStart=/home/multiplane/MultiPlaneDisplay/venv/bin/python /home/multiplane/MultiPlaneDisplay/display.py
Restart=on-failure
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=%t/runtime

[Install]
WantedBy=default.target